# ankavm/oxware-hypervisor.service
[Unit]
Description=OXware Hypervisor Management Service
Documentation=https://github.com/ShinnAsukha/oxware-hypervisor
# network-online.target: ağ gerçekten hazır (sadece yapılandırıldı değil)
# libvirt-guests.service: libvirt hem başladı hem de ağları otomatik açtı
After=network-online.target libvirtd.service libvirt-guests.service
Wants=network-online.target libvirtd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/oxware/oxware
Environment=OXWARE_CONFIG=/etc/oxware/oxware.conf
Environment=PYTHONUNBUFFERED=1

# Dizinleri oluştur
ExecStartPre=/bin/bash -c 'mkdir -p /var/log/oxware /etc/oxware /var/lib/oxware/{isos,disks,backups,templates}'
# libvirtd soketini bekle (reboot sonrası geç hazır olabilir)
ExecStartPre=/bin/bash -c 'for i in $(seq 1 15); do virsh list >/dev/null 2>&1 && break; sleep 2; done; true'
# default ağı başlat (autostart bazen reboot'ta çalışmıyor)
ExecStartPre=/bin/bash -c 'virsh net-list --all 2>/dev/null | grep -q default && virsh net-start default 2>/dev/null || true'
ExecStart=/opt/oxware/venv/bin/python3 /opt/oxware/oxware/backend/app.py
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=mixed
StandardOutput=append:/var/log/oxware/oxware.log
StandardError=append:/var/log/oxware/oxware-error.log
SyslogIdentifier=oxware

# Güvenlik
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
